# resekon/jma-receipt-weborca.service
[Unit]
Description=JMA Receipt WebORCA Service
After=network.target postgresql.service
Requires=postgresql.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=orca
Group=orca
WorkingDirectory=/tmp/weborca
RuntimeDirectory=weborca
RuntimeDirectoryMode=0755
ExecStartPre=/bin/mkdir -p /tmp/weborca
ExecStartPre=/bin/chown orca:orca /tmp/weborca
ExecStartPre=/bin/chmod 755 /tmp/weborca
ExecStart=/usr/bin/jma-receipt-weborca
Restart=on-failure
RestartSec=10
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
